# build/deploy/geth.service
[Unit]
Description=Start geth
Requires=network.target
[Service]
Type=simple
ExecStart=/bin/sh -c "cd /root/blockchain && sh ./scripts/runGeth.sh"
Restart=on-failure
[Install]
WantedBy=default.target
